# etc/blog.conf
[program:blog]
command = /home/zhongxin/Blog/venv/bin/gunicorn -b 128.199.117.54:80 wsgi_gunicorn:app
directory = /home/zhongxin/Blog
user = zhongxin