# etc/blog.conf
[program:blog]
command = /home/zhongxin/.virtualenvs/v_blog/bin/gunicorn -b 128.199.117.54:8080 wsgi_gunicorn:app
directory = /home/zhongxin/Blog
user = zhongxin